# server_openwrt/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
[ -x /etc/hotplug.d/usb/90-usbip ] || chmod 750 /etc/hotplug.d/usb/90-usbip
/etc/hotplug.d/usb/90-usbip 1
exit 0
